# luci/luci-zero-plus/files/etc/uci-defaults/98-uhttpd-default-cgi_prefix.sh
#!/bin/sh
uci batch <<-EOF
	set uhttpd.main.cgi_prefix=/t
	commit uhttpd
EOF
/etc/init.d/uhttpd restart